# conf/springrts-logs-gunicorn.socket
[Unit]
Description=logs.springrts.com socket

[Socket]
ListenStream=/run/gunicorn/infolog.socket

[Install]
WantedBy=sockets.target
